# servers/assets/nginx.conf
start on runlevel [2345]
stop on runlevel [016]
respawn

exec /usr/bin/nginx -c /etc/nginx/conf/nginx.conf  -g "daemon off;"
