# servers/assets/redis.conf
start on runlevel [2345]
stop on runlevel [016]
respawn

exec /etc/redis/bin/redis-server
